# systemd/audit-logs.service
[Unit]
Description=Sweet Narcisse Audit Log Archive Service
Documentation=https://github.com/sweet-narcisse/docs/audit-logs.md
After=network.target

[Service]
Type=oneshot
User=sweet-narcisse
Group=sweet-narcisse
WorkingDirectory=/opt/sweet-narcisse
EnvironmentFile=/opt/sweet-narcisse/.env.production.local
ExecStart=/opt/sweet-narcisse/scripts/archive-audit-logs.sh
StandardOutput=journal
StandardError=journal
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/sweet-narcisse/backups /opt/sweet-narcisse/logs /var/log/sweet-narcisse

# Timeouts
TimeoutStartSec=1800

# Security hardening
NoNewPrivileges=true
ProtectHome=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictRealtime=true

[Install]
WantedBy=multi-user.target
