# systemd/audit-logs.timer
[Unit]
Description=Sweet Narcisse Audit Log Archive Timer
Documentation=https://github.com/sweet-narcisse/docs/audit-logs.md

[Timer]
# Run daily at 5:00 AM
OnCalendar=*-*-* 05:00:00
# Add randomized delay of 0-15 minutes
RandomizedDelaySec=900
# Persist timer across reboots
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target
