# systemd/restore-test.service
[Unit]
Description=Sweet Narcisse Database Restore Test Service
Documentation=https://github.com/sweet-narcisse/docs/backup-guide.md
After=network.target postgresql.service

[Service]
Type=oneshot
User=sweet-narcisse
Group=sweet-narcisse
WorkingDirectory=/opt/sweet-narcisse
EnvironmentFile=/opt/sweet-narcisse/.env.production.local
ExecStart=/opt/sweet-narcisse/scripts/test-restore.sh
StandardOutput=journal
StandardError=journal
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/sweet-narcisse/backups /opt/sweet-narcisse/logs

# Timeouts (restore tests can take a while)
TimeoutStartSec=3600

# Security hardening
NoNewPrivileges=true
ProtectHome=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictRealtime=true

[Install]
WantedBy=multi-user.target
