# systemd/secret-rotation.timer
[Unit]
Description=Sweet Narcisse Secret Rotation Timer
Documentation=https://github.com/sweet-narcisse/docs/secret-rotation.md

[Timer]
# Run on the 1st of each month at 3:00 AM
OnCalendar=*-*-01 03:00:00
# Add randomized delay of 0-30 minutes
RandomizedDelaySec=1800
# Persist timer across reboots
Persistent=true
# Run immediately if last scheduled run was missed
AccuracySec=1h

[Install]
WantedBy=timers.target
